# config
# 1. Create config
mkdir -p /usr/local/etc/reticulum /usr/local/var/log/reticulum
cat > /usr/local/etc/reticulum/config << 'EOF'
enable_transport = Yes

[logging]
loglevel = 3

[interfaces]

[[PostInterface Bridge]]
type = PostInterface
enabled = yes
name = "PostInterface Bridge"
node_url = "https://retichat.com/reticulum"
wake_url = "http://jrl290.ddns.net/v1/wake"
wake_listen_host = 0.0.0.0
wake_listen_port = 4370
poll_interval_seconds = 2.0
bitrate = 1000000000
mtu = 500

[[RMAP]]
type = TCPClientInterface
enabled = yes
target_host = rmap.world
target_port = 4242

[transport]
rns_mtu = 500
pathfinder_max_hops = 128
EOF

# 2. Create rc.d service
cat > /usr/local/etc/rc.d/rnsd << 'EOF'
#!/bin/sh
. /etc/rc.subr
name="rnsd"
rcvar="${name}_enable"
: ${rnsd_enable:="NO"}
: ${rnsd_config:="/usr/local/etc/reticulum"}
command="/usr/local/bin/${name}"
command_args="--config ${rnsd_config} -s"
load_rc_config "${name}"
run_rc_command "$1"
EOF
chmod +x /usr/local/etc/rc.d/rnsd

# 3. Enable and start
sysrc rnsd_enable=YES
service rnsd start

# 4. Check
service rnsd status
tail -f /usr/local/var/log/reticulum/logfile